# deploy/systemd/security-camera-stop.service
[Unit]
Description=Stop security camera headless detector

[Service]
Type=oneshot
ExecStart=/bin/systemctl stop security-camera-headless.service
